# rclip/rclip.service
[Unit]
Description=rclip daemon

[Service]
ExecStart=%h/dotfiles/bin/rclip-daemon

[Install]
WantedBy=default.target
